# systemd/tazgopher@.service
[Unit]
Description=Gopher Per-Connection Server

[Service]
ExecStart=-/opt/tazgopher/bin/tazgopher
StandardInput=socket
